# deploy_tools/gunicorn-upstart.template.conf
description "Gunicorn server for SITENAME"

start on net_device-up
stop on shutdown

respawn

setuid www
chdir /home/www/sites/SITENAME/source

exec ../virtualenv/bin/gunicorn \
    --bind unix:/tmp/SITENAME.socket \
    superlists.wsgi:application
